# udev_scripts/20-opiservo.rules
SUBSYSTEM=="module", DEVPATH=="/module/opiservo", RUN+="/lib/udev/opiservo"
